# etc/inittab.in
::sysinit:/etc/init.d/rcS
::ctrlaltdel:/bin/reboot
::shutdown:/bin/umount -ar -t ext4,vfat
::restart:/init
::respawn:/bin/mdev -df
@extraInittab@
