# etc/inittab
::sysinit:/etc/init.d/rcS
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/init
tty1::respawn:/bin/tinyboot tinyboot.log=debug
tty2::askfirst:/bin/sh
